# init.d.php-fpm.in
#! /bin/sh

### BEGIN INIT INFO
# Provides:          php-fpm2
# Required-Start:    $remote_fs $network
# Required-Stop:     $remote_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts php-fpm2
# Description:       starts the PHP FastCGI Process Manager daemon
### END INIT INFO

prefix=@prefix@
exec_prefix=@exec_prefix@

php_fpm2_BIN=@sbindir@/php-fpm2
php_fpm2_CONF=@sysconfdir@/php-fpm2.conf
php_fpm2_PID=@localstatedir@/run/php-fpm2.pid


php_opts="--fpm-config $php_fpm2_CONF --pid $php_fpm2_PID"


wait_for_pid () {
	try=0

	while test $try -lt 35 ; do

		case "$1" in
			'created')
			if [ -f "$2" ] ; then
				try=''
				break
			fi
			;;

			'removed')
			if [ ! -f "$2" ] ; then
				try=''
				break
			fi
			;;
		esac

		echo -n .
		try=`expr $try + 1`
		sleep 1

	done

}

case "$1" in
	start)
		echo -n "Starting php-fpm2 "

		$php_fpm2_BIN --daemonize $php_opts

		if [ "$?" != 0 ] ; then
			echo " failed"
			exit 1
		fi

		wait_for_pid created $php_fpm2_PID

		if [ -n "$try" ] ; then
			echo " failed"
			exit 1
		else
			echo " done"
		fi
	;;

	stop)
		echo -n "Gracefully shutting down php-fpm2 "

		if [ ! -r $php_fpm2_PID ] ; then
			echo "warning, no pid file found - php-fpm2 is not running ?"
			exit 1
		fi

		kill -QUIT `cat $php_fpm2_PID`

		wait_for_pid removed $php_fpm2_PID

		if [ -n "$try" ] ; then
			echo " failed. Use force-quit"
			exit 1
		else
			echo " done"
		fi
	;;

	status)
		if [ ! -r $php_fpm2_PID ] ; then
			echo "php-fpm2 is stopped"
			exit 0
		fi

		PID=`cat $php_fpm2_PID`
		if ps -p $PID | grep -q $PID; then
			echo "php-fpm2 (pid $PID) is running..."
		else
			echo "php-fpm2 dead but pid file exists"
		fi
	;;

	force-quit)
		echo -n "Terminating php-fpm2 "

		if [ ! -r $php_fpm2_PID ] ; then
			echo "warning, no pid file found - php-fpm2 is not running ?"
			exit 1
		fi

		kill -TERM `cat $php_fpm2_PID`

		wait_for_pid removed $php_fpm2_PID

		if [ -n "$try" ] ; then
			echo " failed"
			exit 1
		else
			echo " done"
		fi
	;;

	restart)
		$0 stop
		$0 start
	;;

	reload)

		echo -n "Reload service php-fpm2 "

		if [ ! -r $php_fpm2_PID ] ; then
			echo "warning, no pid file found - php-fpm2 is not running ?"
			exit 1
		fi

		kill -USR2 `cat $php_fpm2_PID`

		echo " done"
	;;

	configtest)
		$php_fpm2_BIN -t
	;;

	*)
		echo "Usage: $0 {start|stop|force-quit|restart|reload|status|configtest}"
		exit 1
	;;

esac
